# jenkins/jenkins.conf
description "Jenkins"

respawn
start on started network-services
stop on stopping network-services

exec service jenkins start

